# meta-iot-gateway/recipes-network/iotgw-provision/files/iotgw-provision.service
[Unit]
Description=IoT Gateway First-Boot Provisioning
After=local-fs.target
ConditionPathExists=!/var/lib/iotgw-provision.done

[Service]
Type=oneshot
ExecStart=/usr/bin/iotgw-provision.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target

