# meta-iot-gateway/recipes-network/iotgw-provision/files/iotgw-provision.service
[Unit]
Description=IoT Gateway Provisioning (Network + Observability first-boot)
After=local-fs.target
Before=mosquitto.service telegraf.service
RequiresMountsFor=/data

[Service]
Type=oneshot
ExecStart=/usr/bin/iotgw-provision.sh
SyslogIdentifier=iotgw-provision
StandardOutput=journal
StandardError=journal
RemainAfterExit=yes
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_UNIX AF_NETLINK AF_INET AF_INET6
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
# Allow /data mount even when /data/iotgw is absent on fresh systems.
ReadWritePaths=/etc/NetworkManager /etc/default /etc/credstore /etc/mosquitto /var/lib /data

[Install]
WantedBy=multi-user.target
