# meta-iot-gateway/recipes-support/iotgw-pstore-persist/files/iotgw-pstore-persist.service
[Unit]
Description=Prepare persistent systemd-pstore storage on /data
DefaultDependencies=no
After=data.mount overlayfs-setup.service
Before=sysinit.target systemd-pstore.service
Wants=data.mount overlayfs-setup.service
RequiresMountsFor=/data
ConditionPathIsMountPoint=/data

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/sbin/iotgw-pstore-persist
NoNewPrivileges=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=true

[Install]
WantedBy=sysinit.target
