# meta-iot-gateway/recipes-security/iotgw-nftables/files/iotgw-nftables.service
[Unit]
Description=IoT Gateway nftables baseline rules
Documentation=man:nft(8)
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/sbin/nft -f /etc/nftables.conf
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
