# meta-iot-gateway/recipes-security/iotgw-nftables/files/iotgw-nftables.service
[Unit]
Description=IoT Gateway nftables baseline rules
Documentation=man:nft(8)
After=network.target

[Service]
Type=oneshot
# Reset only our baseline table; ignore failure when table is absent.
ExecStartPre=-/usr/sbin/nft delete table inet filter
ExecStart=/usr/sbin/nft -f /etc/nftables.conf
RemainAfterExit=yes
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_NETLINK
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
CapabilityBoundingSet=CAP_NET_ADMIN
AmbientCapabilities=CAP_NET_ADMIN

[Install]
WantedBy=multi-user.target
