# meta-iot-gateway/recipes-support/iotgw-banner/files/iotgw-banner.service
[Unit]
Description=IoT Gateway Dynamic Banner Generator
# Ensure overlay/rauc ready; don't block boot on network-online
After=rauc.service overlayfs-setup.service network.target
Wants=getty.target
Before=getty.target

[Service]
Type=oneshot
ExecStart=/usr/bin/iotgw-banner.sh all
RemainAfterExit=yes
NoNewPrivileges=true
ProtectSystem=true
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6 AF_NETLINK
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true

[Install]
WantedBy=multi-user.target
