# meta-iot-gateway/recipes-support/iotgw-banner/files/iotgw-banner.service
[Unit]
Description=IoT Gateway Dynamic Banner Generator
# Ensure overlay/rauc ready; don't block boot on network-online
After=rauc.service overlayfs-setup.service network.target
Wants=getty.target
Before=getty.target

[Service]
Type=oneshot
ExecStart=/usr/bin/iotgw-banner.sh all
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
