# meta-iot-gateway/recipes-security/iotgw-hardening/files/sshd@.service
[Unit]
Description=OpenSSH Per-Connection Daemon
After=sshdgenkeys.service

[Service]
Environment="SSHD_OPTS="
EnvironmentFile=-/etc/default/ssh
ExecStart=-/usr/sbin/sshd -i $SSHD_OPTS
StandardInput=socket
KillMode=process

# IoT Gateway hardening profile
# Keep privilege escalation available for authorized admin users (sudo/su).
NoNewPrivileges=false
ProtectSystem=full
ProtectHome=false
PrivateTmp=true
PrivateDevices=false
ProtectKernelTunables=false
ProtectKernelModules=false
ProtectKernelLogs=false
ProtectControlGroups=true
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
