# meta-iot-gateway/recipes-security/iotgw-tpm-health/files/iotgw-tpm-health.service
[Unit]
Description=IoT Gateway TPM runtime health probe
After=local-fs.target
RequiresMountsFor=/data
ConditionPathIsMountPoint=/data
ConditionPathExists=/dev/tpmrm0

[Service]
Type=oneshot
ExecStart=/usr/sbin/iotgw-tpm-health
SyslogIdentifier=iotgw-tpm-health
StandardOutput=journal
StandardError=journal
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
ReadWritePaths=-/data/ota/tpm

[Install]
WantedBy=multi-user.target
