# meta-iot-gateway/recipes-support/iotgw-provision/files/iotgw-provision.service
[Unit]
Description=IoT Gateway Provisioning (Network + Observability first-boot)
# NetworkManager is in this recipe's RDEPENDS so it is guaranteed present.
# Ordering after NetworkManager.service avoids the early-boot race where
# nmcli connection reload runs before NM is up and logs
# "failed to reload connections: NetworkManager is not running".
# Deliberately ordering against the unit, NOT network-online.target -- we
# do not want to delay boot waiting for actual link/IP.
After=local-fs.target NetworkManager.service
Wants=NetworkManager.service
Before=mosquitto.service telegraf.service
RequiresMountsFor=/data

[Service]
Type=oneshot
ExecStart=/usr/bin/iotgw-provision.sh
SyslogIdentifier=iotgw-provision
StandardOutput=journal
StandardError=journal
RemainAfterExit=yes
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true
PrivateTmp=true
# fw_setenv needs block-device visibility for /uboot-env backend.
PrivateDevices=false
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_UNIX AF_NETLINK AF_INET AF_INET6
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
# Allow /data mount even when /data/iotgw is absent on fresh systems.
ReadWritePaths=/etc/NetworkManager /etc/default /etc/credstore /etc/mosquitto /var/lib /data

[Install]
WantedBy=multi-user.target
