# meta-iot-gateway/recipes-support/iotgw-provision/files/iotgw-provision.service
[Unit]
Description=IoT Gateway Provisioning (Network + Observability first-boot)
# systemd-networkd is the network stack (RDEPENDS of this recipe). Order after
# it so `networkctl reload` of any provisioned .network unit runs once networkd
# is up. Deliberately ordering against the unit, NOT network-online.target --
# we do not want to delay boot waiting for actual link/IP.
After=local-fs.target systemd-networkd.service
Wants=systemd-networkd.service
Before=mosquitto.service
RequiresMountsFor=/data

[Service]
Type=oneshot
ExecStart=/usr/bin/iotgw-provision.sh
SyslogIdentifier=iotgw-provision
StandardOutput=journal
StandardError=journal
RemainAfterExit=yes
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true
PrivateTmp=true
# fw_setenv needs block-device visibility for /uboot-env backend.
PrivateDevices=false
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_UNIX AF_NETLINK AF_INET AF_INET6
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
# Allow /data mount even when /data/iotgw is absent on fresh systems.
ReadWritePaths=/etc/systemd/network /etc/wpa_supplicant /etc/mosquitto /var/lib /data

[Install]
WantedBy=multi-user.target
