# meta-iot-gateway/recipes-observability/telegraf/files/telegraf.service
[Unit]
Description=Telegraf metrics agent
Documentation=https://github.com/influxdata/telegraf
After=network.target mosquitto.service
Wants=mosquitto.service
ConditionFileNotEmpty=/etc/credstore/telegraf.service.mqtt_username
ConditionFileNotEmpty=/etc/credstore/telegraf.service.mqtt_password
ConditionFileNotEmpty=/etc/credstore/telegraf.service.influxdb_username
ConditionFileNotEmpty=/etc/credstore/telegraf.service.influxdb_password
# Bound restart storms when credentials are missing or invalid.
StartLimitIntervalSec=5min
StartLimitBurst=3

[Service]
Type=simple
User=telegraf
Group=telegraf
StateDirectory=telegraf
UMask=0027
# Non-secret defaults (URL/db name) come from env file.
EnvironmentFile=-/etc/default/iotgw-observability
# Credentials are delivered into $CREDENTIALS_DIRECTORY (a per-unit tmpfs,
# mode 0500) by systemd. The secretstores.systemd Telegraf plugin reads them
# directly — no shell export, no rendered file, no argv exposure.
LoadCredential=mqtt_username:/etc/credstore/telegraf.service.mqtt_username
LoadCredential=mqtt_password:/etc/credstore/telegraf.service.mqtt_password
LoadCredential=influxdb_username:/etc/credstore/telegraf.service.influxdb_username
LoadCredential=influxdb_password:/etc/credstore/telegraf.service.influxdb_password
ExecStart=/usr/bin/telegraf --config /etc/telegraf/telegraf.conf
Restart=on-failure
RestartSec=30
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ProtectClock=yes
ProtectHostname=yes
ProtectControlGroups=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
PrivateTmp=yes
# Keep device namespace isolated for current deployments.
# If Modbus RTU (/dev/ttyUSB*) is enabled in future, this must be revisited
# (disable PrivateDevices or allow explicit device access policy).
PrivateDevices=yes
CapabilityBoundingSet=
AmbientCapabilities=
SystemCallArchitectures=native
SystemCallFilter=@system-service
MemoryDenyWriteExecute=yes
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
RestrictAddressFamilies=AF_UNIX AF_NETLINK AF_INET AF_INET6
LockPersonality=yes
RemoveIPC=yes
ReadWritePaths=/var/lib/telegraf

[Install]
WantedBy=multi-user.target
