# meta-iot-gateway/recipes-ota/ota-updater/files/ota-updater.timer
# SPDX-License-Identifier: MIT
#
# ota-updater.timer: Periodic timer for OTA update checks
#

[Unit]
Description=OTA Update Check Timer
Documentation=man:rauc(1)

[Timer]
# Run 2 minutes after boot, then every hour
OnBootSec=2min
OnUnitActiveSec=1h

# Add randomized delay to avoid thundering herd
RandomizedDelaySec=5min

# Persist timer state across reboots
Persistent=true

[Install]
WantedBy=timers.target
