# meta-iot-gateway/recipes-support/iotgw-machine-id/files/iotgw-machine-id.service
[Unit]
Description=Persist and bind machine-id for immutable rootfs
DefaultDependencies=no
After=data.mount overlayfs-setup.service
# Bind-mount the persistent machine-id over /etc/machine-id BEFORE journald
# starts. journald reads /etc/machine-id once at service start and writes
# its persistent storage to /var/log/journal/<that-id>/. Without this
# ordering it caches the squashfs-baked id, so after the bind-mount lands
# any later `journalctl` (which resolves the current /etc/machine-id)
# looks under the persistent id and finds nothing -- "No journal files
# were found" while journald is in fact writing under the baked id.
Before=sysinit.target systemd-journald.service
Wants=data.mount overlayfs-setup.service
RequiresMountsFor=/data
ConditionPathIsMountPoint=/data

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/sbin/iotgw-machine-id

[Install]
WantedBy=sysinit.target
