# meta-iot-gateway/recipes-support/iotgw-log-persist/files/var-volatile-log-audit.mount
[Unit]
Description=Bind /data-backed audit log storage over /var/volatile/log/audit
Documentation=man:auditd(8)
DefaultDependencies=no
After=data.mount var-volatile.mount iotgw-log-persist-prep.service
Requires=iotgw-log-persist-prep.service
Wants=data.mount
RequiresMountsFor=/data
ConditionPathIsMountPoint=/data
ConditionPathIsDirectory=/data/log/audit
Before=auditd.service umount.target
Conflicts=umount.target

[Mount]
What=/data/log/audit
Where=/var/volatile/log/audit
Type=none
Options=bind,nodev,nosuid,noexec
DirectoryMode=0700
TimeoutSec=10s
